# config/includes.chroot/etc/xdg/autostart/udiskie-automount.desktop
[Desktop Entry]
Name=automounter
Comment=udisk2 automounter with optional notifications;tray icon
Exec=udiskie --tray
